# .config/systemd/user/tray.target
[Unit]
Description=System Tray
Requires=graphical-session-pre.target
